# conf/systemd/affine.service
[Unit]
Description=AFFiNE (YunoHost)
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
User=__APP__
WorkingDirectory=/var/www/__APP__/app
EnvironmentFile=/var/www/__APP__/.env
# Assure l'accès à node/pnpm installés via le helper YunoHost
Environment=PATH=/opt/node_n/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin

# AFFiNE est un monorepo Next.js : on démarre via pnpm (build → start)
ExecStart=/bin/bash -lc 'pnpm start -- --port ${PORT:-4090}'
Restart=always
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target