# conf/systemd/affine.service
[Unit]
Description=AFFiNE (YunoHost smoke test)
After=network.target
Wants=network.target

[Service]
User=__APP__
WorkingDirectory=/var/www/__APP__
Environment=NODE_ENV=production
# Important : PATH inclut le Node géré par YunoHost
Environment=PATH=/opt/node_n/bin:/usr/local/bin:/usr/bin
ExecStart=/bin/bash -lc 'node server.js'
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target